# deploy/systemd/matrix-stock.timer
[Unit]
Description=Run Matrix A-share stock recommendation after market close

[Timer]
OnCalendar=Mon..Fri 20:30:00
Persistent=true
Unit=matrix-stock.service

[Install]
WantedBy=timers.target
