# deploy/systemd/matrix-us.timer
[Unit]
Description=Run Matrix US stock recommendation during China daytime

[Timer]
# 美股白天（中国时区）拉取：此时上一美股交易日已完整收盘，且与晚间 A 股/ETF
# 拉取窗口彻底错开，避免共享 tickflow 免费档 60/min 限额时相互抢占。
OnCalendar=Mon..Fri 14:00:00
Persistent=true
Unit=matrix-us.service

[Install]
WantedBy=timers.target
