# deploy/systemd/matrix-analytics.service
[Unit]
Description=Matrix daily strategy performance analytics
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/Matrix
EnvironmentFile=-/opt/Matrix/.env
ExecStart=/opt/Matrix/scripts/run_analytics.sh
TimeoutStartSec=2h
Nice=10
