# deploy/systemd/matrix-etf.service
[Unit]
Description=Matrix daily ETF recommendation
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/Matrix
EnvironmentFile=-/opt/Matrix/.env
ExecStart=/opt/Matrix/scripts/run_matrix.sh
TimeoutStartSec=4h
Nice=5
