# deploy/systemd/matrix-stock.timer
[Unit]
Description=Run Matrix A-share stock recommendation after market close

[Timer]
OnCalendar=Mon..Fri 19:40:00
Persistent=true
RandomizedDelaySec=5m
Unit=matrix-stock.service

[Install]
WantedBy=timers.target
