# deploy/systemd/matrix-stock.service
[Unit]
Description=Matrix daily A-share stock recommendation
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/Matrix
EnvironmentFile=-/opt/Matrix/.env
ExecStart=/opt/Matrix/scripts/run_stock.sh
TimeoutStartSec=5h
Nice=10
